# sources/supervisor.conf
[program:vms-duo]
command=/opt/vms-duo/vms-duo
autostart=true
autorestart=true
user=vms-duo
stopsignal=KILL 
directory=/opt/vms-duo